# pi-setup/inventar-dashboard.service
[Unit]
Description=Inventar Pro Pi-Dashboard (Web-Admin)
After=network.target inventar-backend.service

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/inventar/backend
Environment=PATH=/home/admin/inventar/backend/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/admin/inventar/backend/.venv/bin/python3 pi_dashboard.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=inventar-dashboard

[Install]
WantedBy=multi-user.target
