# pi-setup/inventar-dashboard.service
[Unit]
Description=Inventar Pro Pi-Dashboard (Web-Admin)
After=network.target inventar-backend.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/inventar/backend
Environment=PATH=/home/admin/inventar/backend/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/admin/inventar/backend/.venv/bin/python3 pi_dashboard.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=inventar-dashboard
ProtectSystem=full
NoNewPrivileges=yes
PrivateTmp=yes
RestrictRealtime=yes
# ProtectHome deaktiviert: Dashboard schreibt audit-log + password-file in ~/

[Install]
WantedBy=multi-user.target
