# pi-setup/inventar-weekly-reboot.service
[Unit]
Description=Inventar Pro Weekly Reboot

[Service]
Type=oneshot
ExecStart=/sbin/shutdown -r +1 "Weekly auto-reboot in 1 minute"
