# pi-setup/inventar-frontend.service
[Unit]
Description=Inventar Pro Frontend (Expo Web)
After=network.target inventar-backend.service

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/inventar/frontend
ExecStart=/usr/bin/npx expo start --web --port 8081 --no-dev
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=CI=1
StandardOutput=journal
StandardError=journal
SyslogIdentifier=inventar-frontend

[Install]
WantedBy=multi-user.target
